# xelatex/install/etc/profile.d/swarmlab-mpi.sh
export PATH=$PATH:/usr/share/swarmlab.io/mpi
